# 80-dvorak.rules
SUBSYSTEM=="input", KERNEL=="event[0-9]*", ACTION=="add", RUN+="/bin/systemctl start dvorak@%k.service"
