# dvorak@.service
[Unit]
Description=Dvorak Virtual Keyboard

[Service]
ExecStart=/usr/bin/dvorak /dev/input/%i keyb k360 k750

[Install]
WantedBy=graphical.target
