# dvorak@.service
[Unit]
Description=Dvorak Virtual Keyboard

[Service]
ExecStart=/usr/local/bin/dvorak -d /dev/input/%i
