# examples/dvorak-usb.service
# /etc/systemd/system/dvorak-usb.service
[Unit]
Description=Dvorak remapping - Glove80 USB
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/local/bin/dvorak-start.sh "MoErgo Glove80 Left Keyboard"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
